# cron/crontab
@hourly echo "Testing cron runner"

@hourly cf run-task cron --command 'echo "Test CF Task"' --name "test-cron-task"

@daily curl -XPOST "http://$GDRIVE_APP_URL:$GDRIVE_APP_PORT/analytics"