# cron/crontab
@hourly echo "Testing cron runner"

@hourly cf run-task cron --command 'echo "Test CF Task"' --name "test-cron-task"

0 12 * * * curl -XPOST "http://$GDRIVE_APP_HOST:$GDRIVE_APP_PORT/analytics"