# cron/crontab
@hourly echo "Hello, world."

@hourly cf run-task cron --command 'echo "Hello, world."' --name "task"
