# cron/crontab
@hourly echo "Testing cron runner"

@hourly cf run-task cron --command 'echo "Test CF Task"' --name "test-cron-task"

*/5 * * * * cf run-task idvametrics --process web
